# .config/systemd/user/alacritty.service
[Unit]
Description=Alacritty Terminal Daemon

[Service]
ExecStart=/usr/bin/alacritty --daemon
Restart=always

[Install]
WantedBy=default.target
